# rootco/containerhost/d0/rootco-opera-web-backup.service
[Unit]
Description=Backup /var/opt/opera/html to cloud.dyn.rootco.de

[Service]
Type=oneshot
ExecStart=/usr/bin/rsync -a /var/opt/opera/html root@cloud.dyn.rootco.de:/backups/{{ salt['grains.get']('id') }}/opera/
