# rootco/containerhost/d0/rootco-minecraft-backup.service
[Unit]
Description=Backup /var/opt/minecraft to cloud.dyn.rootco.de

[Service]
Type=oneshot
ExecStart=/usr/bin/rsync -a /var/opt/minecraft root@cloud.dyn.rootco.de:/backups/{{ salt['grains.get']('id') }}/
