# rootco/containerhost/d0/rootco-znc-backup.service
[Unit]
Description=Backup /var/opt/rootco-znc to cloud.dyn.rootco.de

[Service]
Type=oneshot
ExecStart=/usr/bin/rsync -a /var/opt/rootco-znc root@cloud.dyn.rootco.de:/backups/{{ salt['grains.get']('id') }}/
