# rootco/containerhost/r2d2/certbot.service
[Unit]
Description=Certbot renew
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew
