# rootco/containerhost/r2d2/rootco-measure-download.service
[Unit]
Description=Measure download

[Service]
ExecStart=/usr/local/bin/measure-download.sh
StandardOutput=file:/disk/nextcloud/dlbw.css
