# rootco/containerhost/r2d2/rootco-measure-upload.service
[Unit]
Description=Measure upload

[Service]
ExecStart=/usr/local/bin/measure-upload.sh
StandardOutput=file:/disk/nextcloud/ulbw.css
